# quadlets/open-webui.container
[Unit]
Description=Open WebUI
After=network-online.target ollama-rocm.service
Wants=network-online.target ollama-rocm.service

[Container]
ContainerName=open-webui
Image=ghcr.io/open-webui/open-webui:main
Network=ai-shared.network
PublishPort=3000:8080
Volume=%h/.local/share/open-webui:/app/backend/data:Z
Environment=OLLAMA_BASE_URL=http://ollama:11434

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
