# quadlets/podman-mcp-server.container
[Unit]
Description=Podman MCP Server (HTTP mode)
After=network-online.target podman.socket
Wants=network-online.target podman.socket

[Container]
ContainerName=podman-mcp-server
Image=docker.io/library/node:current-alpine
Network=ai-shared.network
PublishPort=8080:8080
Volume=%t/podman/podman.sock:/run/podman/podman.sock
Environment=CONTAINER_HOST=unix:///run/podman/podman.sock
Exec=npx -y podman-mcp-server@latest --port 8080 --podman-impl=api

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
