# quadlets/vllm-rocm.container
[Unit]
Description=vLLM ROCm OpenAI-compatible API
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=vllm-rocm
Image=docker.io/vllm/vllm-openai-rocm:latest
Network=ai-shared.network
PublishPort=8000:8000
Volume=%h/.cache/huggingface:/root/.cache/huggingface:Z
Volume=%h/.cache/miopen:/workspace/miopen:Z
EnvironmentFile=%h/.config/containers/systemd/stack.env
Environment=MIOPEN_USER_DB_PATH=/workspace/miopen
Environment=MIOPEN_FIND_MODE=FAST
Environment=VLLM_ROCM_USE_AITER=1
Environment=SAFETENSORS_FAST_GPU=1
PodmanArgs=--group-add=video
PodmanArgs=--cap-add=SYS_PTRACE
PodmanArgs=--security-opt=seccomp=unconfined
PodmanArgs=--device=/dev/kfd
PodmanArgs=--device=/dev/dri/renderD128
PodmanArgs=--device=/dev/dri/renderD129
PodmanArgs=--ipc=host
Exec=--model Qwen/Qwen3.5-35B-A3B-FP8 --tensor-parallel 4 -dp 8 --enable-expert-parallel --mm-encoder-tp-mode data --mm-processor-cache-type shm --reasoning-parser qwen3 --enable-prefix-caching --enable-auto-tool-choice --tool-call-parser qwen3_coder

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
